# debian/mlx-openipmi.files
lib/systemd/system/*
var/ipmi_sim/mellanox/sdr.30.main
etc/ipmi/*
usr/bin/*
usr/share/info/OpenIPMI.inf*
usr/share/man/man1/*.1
usr/share/man/man5/*.5
usr/share/man/man7/ipmi_cmdlang.7
usr/share/man/man7/openipmi_conparms.7
usr/share/man/man8/ipmilan.8
/etc/logrotate.d/mlx_ipmid
/etc/logrotate.d/set_emu_param
/etc/rsyslog.d/mlx_ipmid.conf
/etc/rsyslog.d/set_emu_param.conf
